# Yocto/sources/meta-system-controller/recipes-utils/sc-usb-gadget/files/sc-usb-gadget.service
[Unit]
Description=System Controller USB gadget
RequiresMountsFor=/sys/kernel/config

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/usb_gadget_setup.sh
ExecStop=/usr/sbin/usb_gadget_remove.sh

[Install]
WantedBy=sysinit.target
